# server/storage/build_tools/docker/config/supervisor/nginx.conf
[program:nginx]
 command=nginx -g "daemon off;" 
 stderr_logfile=/dz/logs/nginx.err.log 
 stdout_logfile=/dz/logs/nginx.out.log