# server/storage/build_tools/docker/config/supervisor/render.conf
[program:render]
 command=/dz/server/storage/dz_render 
 autostart=true 
 autorestart=true 
 stderr_logfile=/dz/logs/render.err.log 
 stdout_logfile=/dz/logs/render.out.log 