# server/storage/build_tools/docker/config/supervisor/system_helper.conf
[program:system_helper]
 command=/dz/server/storage/dz_system 
 autostart=true 
 autorestart=true 
 stderr_logfile=/dz/logs/system_helper.err.log 
 stdout_logfile=/dz/logs/system_helper.out.log